# daemon/com.pro-gate.review-daemon.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- macOS launchd agent for the pro-review daemon. Rendered + installed by install.sh into
     ~/Library/LaunchAgents/. macOS uses oracle's native browser mode, so there is no
     oracle-chrome counterpart (Chrome is driven per-run). -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.pro-gate.review-daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>@HOME@/.pro-review-daemon/run-daemon.sh</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>10</integer>
  <key>StandardOutPath</key><string>@HOME@/.pro-review-daemon/logs/daemon.out.log</string>
  <key>StandardErrorPath</key><string>@HOME@/.pro-review-daemon/logs/daemon.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key><string>@HOME@</string>
    <key>PATH</key><string>@HOME@/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
