# daemon/pro-review-daemon.service.tmpl
[Unit]
Description=pro-review-daemon: final-tier ChatGPT Pro PR review gate (label: pro-review)
After=network-online.target oracle-chrome.service
Wants=network-online.target
Requires=oracle-chrome.service

[Service]
Type=simple
User=@USER@
Group=@USER@
WorkingDirectory=@HOME@/.pro-review-daemon
ExecStart=@HOME@/.pro-review-daemon/run-daemon.sh
Restart=always
RestartSec=10
KillMode=control-group
TimeoutStopSec=30
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
